# lib/systemd/system/btrfs-snapshots-create.service
[Unit]
Description=Create configured btrfs snapshots

[Service]
Nice=19
IOSchedulingClass=idle
KillSignal=SIGINT
ExecStart=/usr/local/bin/btrfs-snapshots create
