# systemd/btrfs-snapshots-prune.timer
[Unit]
Description=Prune timestamped btrfs snapshots periodically

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
