# systemd/btrfs-snapshots-create.timer
[Unit]
Description=Create timestamped btrfs snapshots periodically

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
