# systemd/btrfs-snapshots-create.timer
[Unit]
Description=create btrfs snapshots every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
